# samba/rootfs/etc/s6-overlay/s6-rc.d/init-smbd/run
#!/command/with-contenv bashio
# ==============================================================================
# Prepare the Samba service for running
# ==============================================================================
declare password
declare username
declare -a interfaces=()
export HOSTNAME

# Check Login data
if ! bashio::config.has_value 'username' || ! bashio::config.has_value 'password'; then
    bashio::exit.nok "Setting a username and password is required!"
fi

bashio::config.require "enabled_shares" "Samba is a tool for sharing folders. Starting it without sharing any folders defeats the purpose."

# Read hostname from API or setting default "hassio"
HOSTNAME=$(bashio::info.hostname)
if bashio::var.is_empty "${HOSTNAME}"; then
    bashio::log.warning "Can't read hostname, using default."
    HOSTNAME="hassio"
fi
bashio::log.info "Hostname: ${HOSTNAME}"

# Determine interfaces list
for interface in $(bashio::config 'interfaces'); do
    if ! bashio::fs.file_exists "/sys/class/net/${interface}/operstate"; then
        bashio::log.warning "Interface ${interface} doesn't exist, skipping it"
    elif ! [[ "$(cat "/sys/class/net/${interface}/operstate")" =~ ^(up|unknown)$ ]]; then
        bashio::log.warning "Interface ${interface} is disabled, skipping it"
    else
        interfaces+=("${interface}")
    fi
done
if [ ${#interfaces[@]} -eq 0 ]; then
    bashio::log.notice "Interfaces: <empty list>"
    bashio::log.notice "Running on localhost only"
    bashio::log.notice "Possible values for the interfaces configuration option are"
    bashio::log.notice "  interfaces:"
    ifconfig | grep -Eo '^\w+' | grep -Ev '^lo|^veth|^docker|^hassio' | while read -r line ; do
        bashio::log.notice "    - ${line}"
    done
else
    bashio::log.info "Interfaces: $(printf '%s ' "${interfaces[@]}")"
fi

# Generate Samba configuration.
jq ".interfaces = $(jq -c -n '$ARGS.positional' --args -- "${interfaces[@]}") |
    .enabled_shares.[] |= ascii_downcase" /data/options.json \
    | tempio \
      -template /usr/share/tempio/smb.gtpl \
      -out /etc/samba/smb.conf

# Init user
username=$(bashio::config 'username')
password=$(bashio::config 'password')
addgroup "${username}"
adduser -D -H -G "${username}" -s /bin/false "${username}"

(echo "$password"; echo "$password") \
    | smbpasswd -a -s -c "/etc/samba/smb.conf" "${username}"
